# supervisord.conf
[supervisord]
nodaemon=true
    
[program:1]
user=root
command=tinyproxy -d -c ./tiny.conf
startsecs = 0
autorestart = true
startretries = 1

